# etc/systemd/system/mydns.service
[Unit]
Description=mydns notify ip

[Service]
Type=simple
ExecStart=/etc/mydns/notify-ip.sh user:password
