# etc/systemd/system/mydns-short.service
[Unit]
Description=mydns notify ip short cycle

[Service]
Type=simple
ExecStart=/etc/mydns/notify-ip-change.sh sample.com
